# src/etc/supervisor.d/php-fpm.ini
[program:php-fpm]
command=/usr/sbin/php-fpm7 -F
autostart=true
autorestart=unexpected
redirect_stderr=True
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
exitcodes=0
