# src/etc/supervisor.d/ttrss.cron.ini
[program:ttrss-cron]
command=/bin/sh -c "while true; do php7 -f /var/www/ttrss/update_daemon2.php ;sleep 5m; done;"
directory=/var/www/ttrss
user=nginx
autostart=true
autorestart=true
redirect_stderr=True
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
exitcodes=0
