# src/etc/supervisor.d/ttrss.cron.ini
[program:ttrss-cron]
command=/bin/sh -c "php -f /var/www/ttrss/update_daemon2.php"
directory=/var/www/ttrss
user=nginx
autostart=true
autorestart=true
redirect_stderr=True
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
exitcodes=0
